# daemon/ops/cliproxyapi.service
[Unit]
Description=CLIProxyAPI for Claude Code
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=linear-daemon
Group=linear-daemon
Environment=HOME=/var/lib/linear-agent-daemon
WorkingDirectory=/var/lib/linear-agent-daemon
ExecStart=/usr/local/bin/cliproxyapi -config /etc/linear-agent-daemon/cliproxyapi.yaml
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
RemoveIPC=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallArchitectures=native
UMask=0077
ReadWritePaths=/var/lib/linear-agent-daemon/.cli-proxy-api

[Install]
WantedBy=multi-user.target
