# daemon/ops/linear-agent-daemon.service
[Unit]
Description=Linear agent webhook daemon
After=network-online.target cliproxyapi.service
Wants=network-online.target cliproxyapi.service

[Service]
Type=simple
User=linear-daemon
Group=linear-daemon
WorkingDirectory=/opt/linear-agent-daemon
EnvironmentFile=/etc/linear-agent-daemon/env
Environment=HOME=/var/lib/linear-agent-daemon
Environment=PATH=/var/lib/linear-agent-daemon/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/node /opt/linear-agent-daemon/dist/index.js
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
# Browser engines (agent screen recordings) build sandboxes from user/pid/net
# namespaces; allow only those kinds rather than none.
RestrictNamespaces=user pid net ipc uts cgroup
RestrictRealtime=true
LockPersonality=true
RemoveIPC=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# @sandbox (systemd >= 255) + unshare/setns let browser engines set up their
# own seccomp/namespace sandboxes; without them they die with SIGSYS at start.
SystemCallFilter=@system-service @sandbox unshare setns
SystemCallArchitectures=native
UMask=0077
ReadWritePaths=/var/lib/linear-agent-daemon

[Install]
WantedBy=multi-user.target
